# supervisord.conf
[program:timer]
command = python /home/dotcloud/current/timetest.py